# greyrhinosec-xray-c2-client/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable xray-client
systemctl restart xray-client
