# greyrhinosec-xray-c2-server/DEBIAN/postinst
#!/bin/bash
set -e
systemctl daemon-reload
systemctl enable xray
systemctl enable nginx
systemctl restart xray
systemctl restart nginx
