# greyrhinosec-xray-c2-client/etc/cron.d/xray-failover
*/5 * * * * root /usr/local/bin/xray-client-failover.sh
